# systemd/pare-daemon.service
# pare-daemon.service — PARE daemon systemd unit (user service).
#
# Operator edits before `systemctl --user enable --now pare-daemon`:
#   - WorkingDirectory: path to the PARE checkout
#   - EnvironmentFile: path to the PARE env file (e.g., /path/to/pare/.env)

[Unit]
Description=Personal Agentic Reverse Engineer — conversational mobile RE operator. (daemon)
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/pare
EnvironmentFile=/path/to/pare/.env
ExecStart=/path/to/pare/.venv/bin/python -m pare
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
