# systemd/{{agent_name}}-daemon.service
[Unit]
Description={{AGENT_DESCRIPTION}} (daemon)
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/{{AGENT_NAME}}
EnvironmentFile=/path/to/{{AGENT_NAME}}/.env
ExecStart=/path/to/{{AGENT_NAME}}/.venv/bin/python -m {{agent_pkg}}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
